# supervisord.conf
[supervisord]


[program:worker]
command=rq worker
process_name=%(program_name)s-%(process_num)s
numprocs=1000
directory=.
stopsignal=TERM
autostart=true
autorestart=true

